# deploy/stfc-scraper.service
[Unit]
Description=STFC Stat Tracker - Daily Pull
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=stfc
Group=stfc
WorkingDirectory=/opt/stfc
ExecStart=/opt/stfc/deploy/run_pull.sh
Environment=HOME=/home/stfc
TimeoutStartSec=600

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stfc-scraper
